# systemd/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for axum-template
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-networks/
After=network-online.target
Wants=network-online.target
BindsTo=axum-template.service
After=axum-template.service

[Service]
Type=simple
User=cloudflared
Group=cloudflared
ExecStart=/opt/axum-template/cf/cloudflared tunnel run axum-template
Restart=on-failure
RestartSec=5
RestartMaxDelaySec=30

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt/axum-template/cf
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
